# deploy/edgeos/watchdog.cron.in
*	*	*	*	*	root	@datadir@/netify_edgeos/watchdog.py
*/2	*	*	*	*	root	@datadir@/watchdog-conf.sh
