# deploy/watchdog.cron
*	*	*	*	*	root	/usr/libexec/netifyd/watchdog.sh
