# deploy/watchdog.cron
*	*	*	*	*	root	/usr/libexec/watchdog.sh
