# deploy/openwrt/files/netifyd.init.in
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2016-2018 eGloo, Incorporated
#
# This is free software, licensed under the GNU General Public License v2.

START=50
STOP=50

USE_PROCD=1
PROG=/usr/sbin/@PACKAGE_TARNAME@
NETIFYD_AUTODETECT=yes

start_service() {
    mkdir -p /var/lib/@PACKAGE_TARNAME@
    mkdir -p /var/run/@PACKAGE_TARNAME@

    source /usr/libexec/@PACKAGE_TARNAME@/functions.sh
    NETIFYD_OPTS=$(auto_detect_options)

    load_modules

    procd_open_instance
    procd_set_param command $PROG -R $NETIFYD_OPTS
    procd_set_param file /etc/@PACKAGE_TARNAME@.conf
    procd_set_param respawn
    procd_close_instance
}
